# system-maintenance.timer
[Unit]
Description=Run system maintenance daily

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true
Unit=system-maintenance.service

[Install]
WantedBy=timers.target
