# deb-package/DEBIAN/postinst
#!/bin/bash

set -e

# Установка почтовой подсистемы, если нет
if ! command -v mail >/dev/null 2>&1; then
    apt update -y
    apt install -y mailutils
fi

# Права на главный скрипт
chmod +x /usr/local/bin/system_maintenance.sh

# Обновляем systemd
systemctl daemon-reload

# Включаем сервис и таймер
systemctl enable system-maintenance.service
systemctl enable --now system-maintenance.timer

exit 0
