# deb-package/DEBIAN/postinst
#!/bin/bash
set -e

echo "[system-maintenance] Installing files..."

# Создание каталога логов
mkdir -p /var/log
touch /var/log/system_maintenance.log
chmod 644 /var/log/system_maintenance.log

# Обновление systemd
systemctl daemon-reload

# Включаем таймер
systemctl enable --now system-maintenance.timer

echo "[system-maintenance] Installation complete."
exit 0
